# com.warp.refresh.fixed.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.warp.refresh</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/usr/local/warp-refresh/warp-refresh.sh</string>
    </array>
    
    <key>StartInterval</key>
    <integer>3600</integer>
    
    <key>RunAtLoad</key>
    <false/>
    
    <key>KeepAlive</key>
    <false/>
    
    <key>StandardOutPath</key>
    <string>/tmp/warp-refresh.out</string>
    
    <key>StandardErrorPath</key>
    <string>/tmp/warp-refresh.err</string>
    
    <key>WorkingDirectory</key>
    <string>/usr/local/warp-refresh</string>
</dict>
</plist>
